# _docker-refactor/crontab
# START CRON JOB
0,15,30,45 * * * * flock /tmp/phoenix.lockfile /app/venv/bin/python3 /app/main.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB
